# crawlers/deploy/amoy-missing.timer
[Unit]
Description=Fill missing blocks at Amoy database

[Timer]
OnBootSec=120s
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
