# crawlers/deploy/arbitrum-one-state.timer
[Unit]
Description=Execute Arbitrum one state crawler each 5m

[Timer]
OnBootSec=15s
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
