# crawlers/deploy/arbitrum-nova-missing.timer
[Unit]
Description=Fill missing blocks at Arbitrum Nova database

[Timer]
OnBootSec=120s
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
