# crawlers/deploy/arbitrum-one-missing.timer
[Unit]
Description=Fill missing blocks at Arbitrum One database

[Timer]
OnBootSec=120s
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
