# crawlers/deploy/avalanche-fuji-missing.timer
[Unit]
Description=Fill missing blocks at Avalanche Fuji database

[Timer]
OnBootSec=120s
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
