# crawlers/deploy/ethereum-historical-crawl-transactions.timer
[Unit]
Description=Runs transactions historical crawler on ethereum

[Timer]
OnBootSec=60s
OnUnitActiveSec=10m

[Install]
WantedBy=timers.target
