# crawlers/deploy/ethereum-metadata.timer
[Unit]
Description=Execute Ethereum metadata crawler each 10m

[Timer]
OnBootSec=20s
OnUnitActiveSec=30m

[Install]
WantedBy=timers.target
