# crawlers/deploy/ethereum-missing.timer
[Unit]
Description=Fill missing blocks at Ethereum database

[Timer]
OnBootSec=40s
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
