# crawlers/deploy/game7-testnet-state.timer
[Unit]
Description=Execute Game7 testnet state crawler each 5m

[Timer]
OnBootSec=15s
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
