# crawlers/deploy/mumbai-missing.timer
[Unit]
Description=Fill missing blocks at Mumbai database

[Timer]
OnBootSec=120s
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
