# crawlers/deploy/mumbai-state-clean.timer
[Unit]
Description=Execute Polygon state clean labels crawler each 25m

[Timer]
OnBootSec=50s
OnUnitActiveSec=25m

[Install]
WantedBy=timers.target
