# crawlers/deploy/polygon-statistics.timer
[Unit]
Description=Update Polygon statistics dashboard each 5 minutes

[Timer]
OnBootSec=10s
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target