# crawlers/deploy/polygon-statistics.timer
[Unit]
Description=Update Polygon statistics dashboards each 6 hours

[Timer]
OnBootSec=10s
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target
