# crawlers/deploy/polygon-metadata.timer
[Unit]
Description=Execute Polygon metadata crawler each 10m

[Timer]
OnBootSec=10s
OnUnitActiveSec=60m

[Install]
WantedBy=timers.target
