# crawlers/deploy/proofofplay-apex-missing.timer
[Unit]
Description=Fill missing blocks at ProofOfPlay database

[Timer]
OnBootSec=120s
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
