# crawlers/deploy/sepolia-state.timer
[Unit]
Description=Execute Sepolia state crawler each 5m

[Timer]
OnBootSec=15s
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
