# crawlers/deploy/wyrm-missing.timer
[Unit]
Description=Fill missing blocks at Wyrm database

[Timer]
OnBootSec=120s
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
