# crawlers/deploy/xdai-missing.timer
[Unit]
Description=Fill missing blocks at XDai database

[Timer]
OnBootSec=35s
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target