# crawlers/deploy/xdai-statistics.timer
[Unit]
Description=Update XDai statistics dashboards each 6 hours

[Timer]
OnBootSec=25s
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target