# crawlers/deploy/zksync-era-testnet-missing.timer
[Unit]
Description=Fill missing blocks at ZkSync Era testnet database

[Timer]
OnBootSec=120s
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target