# crawlers/deploy/zksync-era-testnet-historical-crawl-events.timer
[Unit]
Description=Runs events historical crawler on ZkSync Era testnet

[Timer]
OnBootSec=60s
OnUnitActiveSec=10m

[Install]
WantedBy=timers.target
