# workers/deploy/workers-engine-clean-call-requests.timer
[Unit]
Description=Clean outdated call_requests from Engine DB

[Timer]
OnBootSec=40s
OnUnitActiveSec=60m

[Install]
WantedBy=timers.target